# trading/deploy/systemd/trading-evening.service
[Unit]
Description=Trading LongTerm Evening Decision
After=network-online.target
Wants=network-online.target
OnFailure=trading-unit-failure@%n.service

[Service]
Type=oneshot
WorkingDirectory=/root/.openclaw/workspace/trading
EnvironmentFile=-/root/.openclaw/workspace/trading/.env.longterm
ExecStart=/root/.openclaw/workspace/trading/run_longterm_evening_decision.sh
Nice=10

[Install]
WantedBy=multi-user.target
