# trading/deploy/systemd/trading-evening.timer
[Unit]
Description=Timer for Trading LongTerm Evening Decision

[Timer]
OnCalendar=Mon..Fri 15:35:00
Persistent=true
RandomizedDelaySec=15
Unit=trading-evening.service

[Install]
WantedBy=timers.target
